# bin/docker/8.1/supervisord.conf
[supervisord]
nodaemon=false
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php -d ariables_order=EGPCS /var/www/html/artisan queue:listen --delay=60 --force --memory=2048 --sleep=1 --timeout=3600 --tries=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=8
redirect_stderr=true
stdout_logfile=/var/log/supervisor/supervisord.log
stopwaitsecs=3600